# deploy/systemd/vllm-hust-model@.service.template
[Unit]
Description=vLLM-HUST Model Fleet %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__REPO_DIR__
ExecStart=/usr/bin/env bash __REPO_DIR__/scripts/run_model_fleet_service.sh %i
Restart=always
RestartSec=5
TimeoutStopSec=30
KillMode=process

[Install]
WantedBy=default.target